# config/systemd/user/backup.service
[Unit]
Description=Automated backup.
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/bin/bash /home/%u/bin/backup.sh
